# debian/prerm
#!/bin/bash

set -e

case "$1" in
    remove|upgrade|deconfigure)
        echo "Stopping WildlifeSystems services..."
        
        # Stop services
        for service in /etc/systemd/system/ws-*.service /etc/systemd/system/ws-*.timer; do
            if systemctl is-active "$service" >/dev/null 2>&1; then
                echo "Stopping $service"
                systemctl stop "$service" || true
            fi
            if systemctl is-enabled "$service" >/dev/null 2>&1; then
                echo "Disabling $service"
                systemctl disable "$service" || true
            fi
        done
        
        # Reload systemd
        systemctl daemon-reload || true
        ;;
    
    failed-upgrade)
        # Don't stop services on failed upgrade
        ;;
    
    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
